# forex.service
[Service]
ExecStart=/srv/http/ivasilev/ivasilev --forex
