# website.service.example
[Service]
WorkingDirectory=/srv/http/ivasilev.net
ExecStart=/srv/http/ivasilev.net/ivasilev

[Install]
WantedBy=multi-user.target
