# config/ivasilev.service.example
[Service]
WorkingDirectory=/srv/http/ivasilev
ExecStart=/srv/http/ivasilev/ivasilev

[Unit]
Requires=mongodb.service
After=mongodb.service
