# config/ivasilev.forex.service.example
[Service]
WorkingDirectory=/srv/http/ivasilev
ExecStart=/srv/http/ivasilev/ivasilev.forex

[Unit]
Requires=mongodb.service
After=mongodb.service
